# install/com.lanfirst.resolverd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Template: install.sh substitutes __BIN__ and __HOME__ before installing this
     to /Library/LaunchDaemons/com.lanfirst.resolverd.plist (root-owned). This is
     the only root component: it reconciles /etc/resolver to the user's config so
     no per-domain sudo is ever needed. See docs/adr/0003. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lanfirst.resolverd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BIN__</string>
        <string>-config</string>
        <string>__HOME__/.config/lanfirst/config.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/lanfirst-resolverd.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/lanfirst-resolverd.log</string>
</dict>
</plist>
